# roles/mysqlserver/files/autoxtrabackup.cron
# cron.d task for autoextrabackup
# m h dom mon dow user	command
0 *   * * *   root    /srv/scripts/autoxtrabackup.sh
